# .config/systemd/user/monitor-hotplug.service
[Unit]
Description=Monitor hotplug listener

[Service]
ExecStart=%h/.local/bin/screen_monitor
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
